# infra/reconciler/reconciler.timer
[Unit]
Description=Run cella reconciler every 20s
Documentation=https://github.com/cellajs/cella info/DEPLOYMENT_PLATFORM_DESIGN.md

[Timer]
# Fire 30s after boot so docker + cloud-init have settled, then every 20s.
# 20s is the cutover-latency budget: CI writes the tag, then waits up to
# ~5 min for /health to flip. With 20s ticks, expected median observed
# delay is ~10s plus pull + restart time (~20–40s for our images).
OnBootSec=30s
OnUnitActiveSec=20s
# Catch up if the system was suspended/asleep — we don't want a host that
# was offline during a deploy to stay on the old tag indefinitely after wake.
Persistent=true
AccuracySec=2s

[Install]
WantedBy=timers.target
